# etc/cron.d/update-index-txt.crontab
# Update index.txt periodically.
01 * * * * root sleep $(($$ \% 900)); @target@/update-index-txt
