# etc/cron.d/update-crl.crontab
# Update CRL periodically.
# How often depends on "default_crl_days" parameter in openssl.cnf.
0 4 * * 7 root sleep $(($$ \% 900)); @target@/easy-rsa/update-crl
